# docker/app/rootfs/etc/cron.d/scheduler
* * * * * php /app/artisan schedule:run -n | grep -v 'No scheduled commands are ready to run.'
